# php-fpm/context/etc/profile.d/python.sh
# set alias python to python3 if available
if [ -f '/usr/bin/python3' ]; then
  alias python='python3'
fi